# packaging/nctv-player_1.0.2_armhf/DEBIAN/postinst
#!/bin/bash
# postinst - Post-installation script for nctv-player .deb package
set -e

echo "[nctv-player] Running post-installation setup..."

# Create required directories
mkdir -p /var/lib/nctv-player/playlist/playlist-background
mkdir -p /var/lib/nctv-player/playlist/playlist-main
mkdir -p /var/lib/nctv-player/playlist/playlist-horizontal
mkdir -p /var/lib/nctv-player/playlist/playlist-vertical
mkdir -p /etc/nctv-player
mkdir -p /var/log

# Create default config if it doesn't exist
CONFIG_FILE="/etc/nctv-player/config.ini"
if [ ! -f "$CONFIG_FILE" ]; then
    echo "[nctv-player] Creating default configuration at $CONFIG_FILE"
    cat > "$CONFIG_FILE" <<EOF
[General]
kioskMode=true
retryIntervalMs=5000
imageDurationMs=10000
audioEnabled=false

[Paths]
playlistRoot=/var/lib/nctv-player/playlist
logPath=/var/log/nctv-player.log

[Display]
targetWidth=1920
targetHeight=1080

[Optimization]
optimizedSuffix=_optimized
EOF
fi

# Set permissions
chmod 755 /usr/bin/nctv-player
chown -R root:root /var/lib/nctv-player
chmod -R 755 /var/lib/nctv-player

# Install and enable systemd service
if [ -f /lib/systemd/system/nctv-player.service ]; then
    echo "[nctv-player] Enabling systemd service..."
    systemctl daemon-reload
    systemctl enable nctv-player.service
    systemctl restart nctv-player.service || true
fi

echo "[nctv-player] Installation complete."
exit 0
